# supervisord.conf
; Supervisor program group for Epicenter Mail.
; Copied to /etc/supervisor/conf.d/mailing_epicenter.conf
; Program names are prefixed "epimail-" so they never collide with the
; existing "nexomail-" group running on this host. Ports 8002 (serve) and
; 8082 (reverb) are likewise distinct from the nexomail group (8000/8080).

[program:epimail-serve]
process_name=%(program_name)s
command=/usr/bin/php8.4 /var/www/mailing.epicentersport.com/backend/artisan serve --host=127.0.0.1 --port=8002
directory=/var/www/mailing.epicentersport.com/backend
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/mailing.epicentersport.com/backend/storage/logs/serve.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=10

[program:epimail-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php8.4 /var/www/mailing.epicentersport.com/backend/artisan queue:work redis --tries=3 --max-time=3600
directory=/var/www/mailing.epicentersport.com/backend
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/mailing.epicentersport.com/backend/storage/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

[program:epimail-reverb]
process_name=%(program_name)s
command=/usr/bin/php8.4 /var/www/mailing.epicentersport.com/backend/artisan reverb:start --host=127.0.0.1 --port=8090
directory=/var/www/mailing.epicentersport.com/backend
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/mailing.epicentersport.com/backend/storage/logs/reverb.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=10

[program:epimail-scheduler]
process_name=%(program_name)s
command=/usr/bin/php8.4 /var/www/mailing.epicentersport.com/backend/artisan schedule:work
directory=/var/www/mailing.epicentersport.com/backend
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/mailing.epicentersport.com/backend/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=60
