# supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:laravel-serve]
command=php artisan serve --host=127.0.0.1 --port=8000
directory=/var/www/nexosmart-mail/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/laravel-serve.err.log
stdout_logfile=/var/log/laravel-serve.out.log
user=www-data

[program:laravel-queue]
command=php artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
directory=/var/www/nexosmart-mail/backend
autostart=true
autorestart=true
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/laravel-queue.err.log
stdout_logfile=/var/log/laravel-queue.out.log
user=www-data
stopwaitsecs=3600

[program:laravel-reverb]
command=php artisan reverb:start --host=127.0.0.1 --port=8080
directory=/var/www/nexosmart-mail/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/laravel-reverb.err.log
stdout_logfile=/var/log/laravel-reverb.out.log
user=www-data

[program:laravel-scheduler]
command=sh -c "while true; do php artisan schedule:run --verbose --no-interaction >> /var/log/laravel-scheduler.log 2>&1; sleep 60; done"
directory=/var/www/nexosmart-mail/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/laravel-scheduler.err.log
stdout_logfile=/var/log/laravel-scheduler.out.log
user=www-data
